# united_docker/supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:unit3d-queue]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan queue:work
startsecs = 0
autostart=true
autorestart=true
user=www-data
numprocs=10